# supervisor.conf/rtmp_rtmp.ini
[program:rtmp_rtmp.ini]
command=ffmpeg -i "rtmp://14.161.5.88:1935/live/liveStream" -vcodec copy -acodec aac -strict -2 -ar 48000 -b:a 128k -f flv rtmp://118.69.132.214:1935/Lasta
redirect stderr=true
stdout_logfile= /var/log/supervisor/rtmp_rtmp.ini.log
